# defaults/.config/systemd/user/sync-mail.service
[Unit]
Description=Sync mail

[Service]
Type=oneshot
ExecStart=%h/.scripts/mail/sync-mail.sh
